# .config/systemd/user/xdg-desktop-portal-gnome.service.d/wait-for-niri.conf
# Durable fix for screen sharing (Zoom/OBS/etc.) failing on niri.
#
# Symptom: the ScreenCast portal advertised AvailableSourceTypes=0 (no monitor,
# window, or virtual sources), so clients had nothing to share.
#
# Cause: a startup race. niri provides screencast via the org.gnome.Mutter.ScreenCast
# D-Bus name, which it registers ASYNCHRONOUSLY after its process starts. The vendor
# unit only orders the portal After=graphical-session.target, so xdg-desktop-portal-gnome
# could activate and query Mutter.ScreenCast before that name existed, then latch a
# permanently-degraded "0 sources" state until manually restarted.
#
# Fix: order after niri's compositor service and block portal startup until the
# Mutter.ScreenCast name is actually claimable on the session bus.

[Unit]
After=wayland-wm@niri.desktop.service

[Service]
# Wait up to 30s for niri to publish the interface before the portal starts.
# Leading '-' means: if it ever times out, start anyway (no worse than the
# vendor behavior). gdbus ships with glib2.
ExecStartPre=-/usr/bin/gdbus wait --session --timeout 30 org.gnome.Mutter.ScreenCast
