# dot/profile.d/40_dircolors
#!/bin/sh

# Dir Colors
if which dircolors >/dev/null 2>&1; then
	eval $(dircolors)
	export LS_COLORS
fi
