# dot/profile.d/50_ssh_askpass
#!/bin/sh

# SSH
if [ -e /usr/lib/openssh/gnome-ssh-askpass ]; then
	export SSH_ASKPASS=/usr/lib/openssh/gnome-ssh-askpass
fi
