# dot/profile.d/99_xinit
#!/bin/sh

if [ -z "$DISPLAY" ]; then
    case "$(tty)" in
        /dev/tty1)
            exec xinit
            ;;
    esac
fi
